# conf/supervisord.ini
[program:webvirtcloud]
command=/usr/local/bin/gunicorn webvirtcloud.wsgi:application -c /opt/webvirtcloud/gunicorn.conf.py
directory=/opt/webvirtcloud
user=root
autostart=true
autorestart=true
redirect_stderr=true

[program:novncd]
command=/usr/bin/python3 /opt/webvirtcloud/console/novncd
directory=/opt/webvirtcloud
user=root
autostart=true
autorestart=true
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
